# TODO/stage6/01-tweaks/includes.chroot/etc/skel/.config/autostart/keymap.desktop
[Desktop Entry]
Name=KeyMap
Exec=/bin/keymap
Terminal=false
Type=Application
StartupNotify=false